# telegram-claude-bot.service
[Unit]
Description=Telegram Claude Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=criznguyen
WorkingDirectory=/home/criznguyen/projects/telegram-claude-bot
EnvironmentFile=/home/criznguyen/projects/telegram-claude-bot/.env
ExecStart=/home/criznguyen/.pyenv/versions/3.11.9/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
